# my_flask_app.service
[Unit]
Description=My Flask Web App
After=network.target

[Service]
User=azureuser # Ensure this is correct (azureuser on Azure VM, navabshariff on local)
Group=azureuser # Ensure this is correct

# The absolute path to your application's directory
WorkingDirectory=/opt/my_flask_app

# Load environment variables from a file that will be created by custom_data
EnvironmentFile=/etc/my_flask_app.env

# The command to execute: Call the wrapper script
ExecStart=/bin/bash /opt/my_flask_app/start_app_in_venv.sh

Restart=on-failure
RestartSec=10

StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target